# gentoo/portage/net-ftn/qicoxe/files/qico.initd
#!/sbin/runscript

start() {
	ebegin "Starting qico daemon"
	start-stop-daemon --start -c fido:fido --pidfile /var/run/qico/qico.pid --exec /usr/sbin/qico -- -d
	eend $?
}

stop() {
	ebegin "Stopping qico daemon"
	start-stop-daemon --stop --pidfile /var/run/qico/qico.pid
	eend $?
}
